# deploy/horizon.conf
[program:]
process_name=%(program_name)s
command=
autostart=true
autorestart=true
user=
redirect_stderr=true
stdout_logfile=
stopwaitsecs=3600